# winclip.service
[Unit]
Description=WinClip clipboard daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%WINCLIP_DIR%
ExecStart=%WINCLIP_DAEMON_CMD%
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target